# distnoted_watchdog/pkgbuild/scripts/postinstall
#!/bin/bash

/bin/launchctl load "/Library/LaunchDaemons/uk.co.amsys.distnoted_watchdog.plist"

exit 0
